# infra/ansible/roles/postgresql/templates/pg-backups.timer
[Unit]
Description=Postgresql Backup

[Timer]
OnCalendar={{postgresql_backup_schedule}}

[Install]
WantedBy=timers.target
