# infra/ansible/roles/django/templates/ml-train.timer.j2
[Unit]
Description=Train ML

[Timer]
OnCalendar={{ml_train_schedule}}

[Install]
WantedBy=timers.target
