# infra/ansible/roles/django/templates/tournesol-yt-playlists.timer.j2
[Unit]
Description=Tournesol YT playlists schedule

[Timer]
OnCalendar={{tournesol_yt_playlists_schedule | default('@0') }}
Persistent=yes

[Install]
WantedBy=timers.target
