# infra/ansible/roles/django/templates/tournesol-api-registration-reminder.timer.j2
[Unit]
Description=Tournesol API registration reminder

[Timer]
OnCalendar={{tournesol_api_registrationreminder_schedule}}
Persistent=yes

[Install]
WantedBy=timers.target
