# infra/ansible/roles/django/templates/tournesol-api-remind-registrations.timer.j2
[Unit]
Description=Tournesol API remind registrations

[Timer]
OnCalendar={{tournesol_api_remindregistrations_schedule}}
Persistent=yes

[Install]
WantedBy=timers.target
