# infra/ansible/roles/django/templates/tournesol-api-delete-inactive-users.timer.j2
[Unit]
Description=Tournesol API inactive users clean-up

[Timer]
OnCalendar={{tournesol_api_deleteinactiveusers_schedule}}
Persistent=yes

[Install]
WantedBy=timers.target
