# infra/ansible/roles/django/templates/tournesol-twitterbot.timer.j2
[Unit]
Description=Tournesol twitterbot schedule

[Timer]
OnCalendar={{tournesol_twitterbot_schedule | default('') }}
Persistent=yes

[Install]
WantedBy=timers.target
