# infra/ansible/roles/django/templates/tournesol-api-cleartokens.timer.j2
[Unit]
Description=Tournesol API tokens clean-up

[Timer]
OnCalendar={{tournesol_api_cleartokens_schedule}}
Persistent=yes

[Install]
WantedBy=timers.target
