# infra/ansible/roles/django/templates/tournesol-api-create-dataset.timer.j2
[Unit]
Description=Tournesol API create dataset

[Timer]
OnCalendar={{tournesol_api_createdataset_schedule}}

[Install]
WantedBy=timers.target
