# infra/ansible/roles/backups/templates/export-backups.timer.j2
[Unit]
Description=Export Backups to S3

[Timer]
OnCalendar={{export_backups_schedule}}

[Install]
WantedBy=timers.target
