# infra/ansible/roles/mediawiki/templates/mediawiki-backups.timer
[Unit]
Description=Mediawiki XML Backup

[Timer]
OnCalendar={{mediawiki_backup_schedule}}

[Install]
WantedBy=timers.target
