# infra/ansible/roles/monitoring/files/external-urls-monitoring.service
[Unit]
Description=External URLs monitoring

[Service]
Type=oneshot
ExecStart=/usr/local/bin/external-urls-monitoring.sh
DynamicUser=yes
