# infra/ansible/roles/monitoring/files/external-urls-monitoring.timer
[Unit]
Description=External URLs monitoring

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
